# files/local-ports.cron
*/5 * * * * root /usr/local/bin/local-ports.sh
